# todo.conf
[program:todo]
command=/python3 server.py
directory=/var/todo
autostart=true
autorestart=true
